# eink-art-gallery/rootfs/etc/services.d/eink-art-gallery/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
